# 95-yubikey-hotplug.rules
ACTION=="remove", ENV{ID_VENDOR_ID}=="1050", ENV{ID_MODEL_ID}=="0116", RUN+="/usr/bin/yubikey-hotplug disconnect"
